# systemd/darkwg-api.service
[Unit]
Description=DarkWG API
After=network.target darkwg-quick@darkwg0.service
Requires=darkwg-quick@darkwg0.service

[Service]
Type=simple
WorkingDirectory=/opt/darkwg
EnvironmentFile=/etc/darkwg/api.env
ExecStart=/opt/darkwg/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port __API_PORT__
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
